# ops/project-giftboard-frontend@.service
[Unit]
Description=Project Frontend (%i)
Documentation=Deploys static frontend to /srv/projects/%i/frontend/current
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/sh -lc 'test -d /srv/projects/%i/frontend/current'
ExecStartPost=/bin/sh -lc 'if command -v nginx >/dev/null 2>&1; then nginx -t && systemctl reload nginx; fi'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
